# vps/notebooklm-research.timer
[Unit]
Description=Lancia la ricerca NotebookLM ogni mattina

[Timer]
# Ora locale del VPS (fuso impostato a Europe/Rome). 07:00 ora italiana.
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
