# vps/notebooklm-podcast.timer
[Unit]
Description=Lancia KANRI Audio ogni lunedì mattina

[Timer]
# Ora locale del VPS (fuso Europe/Rome). Lunedì alle 08:00.
OnCalendar=Mon *-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
