# vps/notebooklm-watcher.timer
[Unit]
Description=Controlla Notion ogni 5 minuti per le news spuntate

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
